# jd_telegram-cli/lib/jd_ShareCode.service
[Unit]
Description='send jd share code to JDShareCodebot"

[Service]
ExecStart=/bin/bash /root/script/jd_ShareCode.sh
