# jd_telegram-cli/lib/jd_ShareCode.timer
[Unit]
Description="run jd_ShareCode.service every Monday"

[Timer]
OnCalendar=Mon, 00:01

[Install]
WantedBy=multi-user.target
